# deploy/systemd/youth-policy-api.service.example
[Unit]
Description=Youth Policy Chatbot FastAPI
After=network.target

[Service]
Type=simple
User=youth-policy
Group=youth-policy
WorkingDirectory=/opt/youth-policy-rag
EnvironmentFile=/etc/youth-policy-chatbot/api.env
ExecStart=/opt/youth-policy-rag/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
